# crons/crontab
# m h  dom mon dow   command
*/5 * * * * /usr/local/hfa/crons/update_scoreboard.sh
*/20 * * * * /usr/local/hfa/crons/get_checkins.py && /usr/local/hfa/crons/compute_stats.py
